# contrib/systemd/moximo-master.service
[Unit]
Description=Moximo Master

[Service]
ExecStart=/usr/bin/moximo-master

[Install]
WantedBy=multi-user.target
